# whatsapp_bot/soonhoe-bot.service
[Unit]
Description=Soon Hoe WhatsApp Bot
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/soonhoe-bot
EnvironmentFile=/opt/soonhoe-bot/.env
ExecStart=/opt/soonhoe-bot/venv/bin/gunicorn \
    --workers 2 \
    --bind 0.0.0.0:8080 \
    --timeout 60 \
    --access-logfile logs/access.log \
    --error-logfile logs/error.log \
    wsgi:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
